# etc/supervisord.conf
[program:backup_worker]
command={{ PYTHON }} {{ PROJECT_DIR }}/manage.py rqworker default

[program:bak_server]
command={{ PYTHON }} {{ PROJECT_DIR }}/manage.py runserver 0.0.0.0:8080
